# nibelite-app-subscriptions/scripts/smsnews.cron
#!/bin/sh
15 11 * * * root /opt/nibelite/apps/smsnews-notifier
